# server/systemd/zzx-bitnodes-crawler.service
[Unit]
Description=ZZX-Labs Bitnodes Daemon
Documentation=https://zzx-labs.io/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=zzx
Group=zzx

WorkingDirectory=/srv/zzx-labs

Environment=PYTHONUNBUFFERED=1
Environment=ZZX_LABS_ROOT=/srv/zzx-labs
Environment=ZZX_LABS_RUN_DIR=/srv/zzx-labs/run
Environment=ZZX_LABS_LOG_DIR=/srv/zzx-labs/log
Environment=ZZX_LABS_API_DIR=/srv/zzx-labs/bitcoin/bitnodes/api

EnvironmentFile=-/etc/zzx-labs/bitnodes.env

ExecStart=/usr/bin/python3 /srv/zzx-labs/tools/bitnodes/bitnodes.py daemon
ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=15

TimeoutStopSec=120

NoNewPrivileges=true
PrivateTmp=true

StandardOutput=append:/srv/zzx-labs/log/bitnodesd.log
StandardError=append:/srv/zzx-labs/log/bitnodesd.error.log

[Install]
WantedBy=multi-user.target
